# simulatedvehicle/canreplay/service/evcansimulation.blog-vehicle-01.service
[Unit]
Description=CAN Simulator

[Service]
Restart=always
User=ubuntu
RestartSec=1
ExecStart=/bin/bash /home/ubuntu/aws-iot-fleetwise-evbatterymonitoring/simulatedvehicle/canreplay/service/start_simulation.sh blog-vehicle-01
WorkingDirectory=/home/ubuntu/aws-iot-fleetwise-evbatterymonitoring/simulatedvehicle/canreplay/service

[Install]
WantedBy=multi-user.target
